# .config/systemd/.config/systemd/user/xremap.service
[Unit]
Description=xremap key remapping daemon

[Service]
ExecStart=/usr/bin/xremap /home/xiiiihao/.config/xremap/config.yml
Restart=always

[Install]
WantedBy=default.target
